# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/study_foolog/bin/uwsgi --ini /srv/study_foolog/.config/uwsgi/uwsgi-app.ini